# kinship-backend.service
[Unit]
Description=Amarktai Network Backend API
After=network.target mongod.service
Wants=mongod.service

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/amarktai/backend
EnvironmentFile=/var/www/amarktai/backend/.env
ExecStart=/var/www/amarktai/venv/bin/uvicorn server:app --host 127.0.0.1 --port 8000 --workers 2
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=amarktai-backend

[Install]
WantedBy=multi-user.target
